# 036-daemon-monitor/daemon36.service
[Unit]
Description=Sample daemon process

[Service]
Type=simple
ExecStart=/mydaemon/bin/daemon36

[Install]
WantedBy=multi-user.target